# config/etc/crontab
RUBY_MANAGER=chruby
RAILS_ENV=production
MAILTO=webmaster@mushroomobserver.org

*   * * * * /var/web/mo/bin/run script/update_ip_stats.rb
# retransfer_images and the find -mmin +3 | rsync --remove-source-files
# line below it are both retired (see GitHub issue #4791): the rsync
# raced Image::Processor#process -- it could delete a locally-complete
# original before the derived sizes were generated, since it had no idea
# any of that was still in flight. TransferImagesJob (config/
# recurring.yml doesn't schedule it -- it's event-driven, started when
# processing/an iNat import batch finishes) replaces the transfer side;
# StaleImageFilesJob (every hour) and ImageGapDetectorJob (weekly) are
# the safety nets -- see #4791's target design for the 4-part pipeline.
*/4 * * * * /var/web/mo/bin/run script/parse_log
13  3 * * * /var/web/mo/bin/run script/backup_blocked_ips.sh
# This has been superseded by TransferImagesJob/StaleImageFilesJob/
# ImageGapDetectorJob (see above):
# 13  4 * * * /var/web/mo/bin/run script/verify_images
43  4 * * * /var/web/mo/bin/run script/refresh_sitemap
43  5 * * * /var/web/mo/bin/run script/backup_database.sh
13  6 * * * /var/web/mo/bin/run script/update_googlebots.rb
43  6 * * * /var/web/mo/bin/run script/update_snapshots
13  9 * * * sudo /usr/bin/certbot renew --quiet
